# runit/nginx/run
#!/bin/sh
exec chpst -u www-data /usr/sbin/nginx -c /rails-app/nginx.conf
