# runit/nginx/run
#!/bin/sh
exec chpst -unginx /usr/sbin/nginx -c /rails-app/nginx.conf