# runit/rails/run
#!/bin/sh
cd /rails-app
exec /usr/bin/ruby /rails-app/bin/bundle exec /usr/bin/unicorn -o 0.0.0.0 -p 8080 --no-default-middleware