# runit/rails/run
#!/bin/sh
cd /rails-app
export RAILS_ENV=production
exec chpst -uwww-data /usr/local/bin/ruby /rails-app/bin/bundle exec unicorn -o 0.0.0.0 -p 8080 --no-default-middleware
